# bc-zfsonlinux-memorymanagement2.crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root bc-zfsonlinux-memorymanagement2.py >/dev/null 2>&1
